# crates/domains/vaultwarden/templates/vaultwarden-backup.timer
[Unit]
Description=Daily Vaultwarden Backup

[Timer]
OnCalendar=daily
RandomizedDelaySec=3600
Persistent=true

[Install]
WantedBy=timers.target
